# recipes-core/tpm2-util/files/mender-luks-tpm-key-watcher.path
[Path]
PathChanged=@@MENDER/LUKS_KEY_FILE@@
MakeDirectory=True

[Install]
WantedBy=multi-user.target
