# recipes-core/tpm2-util/files/mender-luks-tpm-key-watcher.path
[Unit]
After=data.mount
Requires=data.mount

[Path]
PathChanged=@@MENDER/LUKS_LEGACY_KEY_FILE@@
MakeDirectory=True

[Install]
WantedBy=multi-user.target
